# etc/crontab/tail
    # Daily forecast:
@@(C:MIT_CHAIN_START_HOUR_CRONTAB_TZ) * * * @@(C:MIT_SCRDIR)/mit_is_login_node.ksh && @@(C:MIT_SCRDIR)/mit_start.ksh --pass --job-multiple --without-kill --rundate $(date +\%Y\%m\%d) --force --try-resume 1> @@(C:MIT_VARDIR)/log/$(hostname).$(date +\%Y\%m\%d).mit_start.out 2>&1
@@(C:MIT_CHAIN_CHECK_HOUR_CRONTAB_TZ) * * * @@(C:MIT_SCRDIR)/mit_is_login_node.ksh && @@(C:MIT_SCRDIR)/mit_check_run.ksh --refresh --notify --weekday @@(C:MIT_WEEKDAY_NAME) --run 1> @@(C:MIT_VARDIR)/log/$(hostname).$(date +\%Y\%m\%d).mit_check_run.out 2>&1
   
   # Reset crontab:
00 00 * * * @@(C:MIT_SCRDIR)/mit_setup_crontab.ksh --force 1> @@(C:MIT_VARDIR)/log/$(hostname).$(date +\%Y\%m\%d).mit_setup_crontab.out 2>&1
