# etc/crontab/weekday
   # @@(C:OPA_SUBMIT_WEEKDAY_NAME) run (opa_weekday is @@(C:OPA_WEEKDAY_NAME))
@@(C:OPA_CRON_COMMENT)@@(C:OPA_DOWNLOAD_START_HOUR_CRONTAB_TZ) * * @@(C:OPA_SUBMIT_WEEKDAY) @@(C:OPA_SCRDIR)/opa_is_login_node.ksh && @@(C:OPA_BINDIR)/timeout 3600 @@(C:OPA_SCRDIR)/opa_get.ksh --weekday @@(C:OPA_WEEKDAY_NAME) --force-download --without-timeout 1> @@(C:OPA_VARDIR)/log/$(hostname).$(date +\%Y\%m\%d).opa_get.out 2>&1
@@(C:OPA_CHAIN_START_HOUR_CRONTAB_TZ) * * @@(C:OPA_SUBMIT_WEEKDAY) @@(C:OPA_SCRDIR)/opa_is_login_node.ksh && @@(C:OPA_SCRDIR)/opa_start.ksh --pass --job-multiple --force --with-timeout --weekday @@(C:OPA_WEEKDAY_NAME) --try-resume 1> @@(C:OPA_VARDIR)/log/$(hostname).$(date +\%Y\%m\%d).opa_start.out 2>&1
@@(C:OPA_CRON_COMMENT)@@(C:OPA_CHAIN_CHECK_HOUR_CRONTAB_TZ) * * @@(C:OPA_CHECK_WEEKDAY) @@(C:OPA_SCRDIR)/opa_is_login_node.ksh && @@(C:OPA_SCRDIR)/opa_check_run.ksh --refresh --notify --weekday @@(C:OPA_WEEKDAY_NAME) 1> @@(C:OPA_VARDIR)/log/$(hostname).$(date +\%Y\%m\%d).opa_check_run.out 2>&1
0 0 * * @@(C:OPA_CHECK_WEEKDAY) @@(C:OPA_SCRDIR)/opa_is_login_node.ksh && @@(C:OPA_SCRDIR)/opa_cleaner.ksh --force --weekday @@(C:OPA_WEEKDAY_NAME) 1> @@(C:OPA_VARDIR)/log/$(hostname).$(date +\%Y\%m\%d).opa_cleaner.out 2>&1
